# com.mattwisse.autopost-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <!-- Unique label for launchd -->
  <key>Label</key>
  <string>com.mattwisse.autopost-watcher</string>

  <!-- The script to run -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/mattwisse/Documents/Claude/Projects/Youtube Channels Project/auto_watcher.sh</string>
  </array>

  <!-- Start automatically when you log in -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Restart if it ever crashes -->
  <key>KeepAlive</key>
  <true/>

  <!-- Log output here -->
  <key>StandardOutPath</key>
  <string>/Users/mattwisse/Documents/Claude/Projects/Youtube Channels Project/watcher.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mattwisse/Documents/Claude/Projects/Youtube Channels Project/watcher_error.log</string>

  <!-- Working directory -->
  <key>WorkingDirectory</key>
  <string>/Users/mattwisse/Documents/Claude/Projects/Youtube Channels Project</string>

</dict>
</plist>
